# vpndnsutils/DEBIAN/postinst
systemctl enable vpndnsutils
systemctl start vpndnsutils